# UBUNTU-RCONFIG INSTALL.conf
######################################################################################
cd ~/
sudo apt -y install wget
wget https://www.rconfig.com/downloads/rconfig6_ubuntu.sh -O install.sh
chmod +x install.sh
apt update && apt upgrade -y
sudo ./install.sh
######################################################################################
mysql -u root -p

CREATE DATABASE rconfig;
CREATE USER 'rconfig'@'localhost' IDENTIFIED BY 'BKqD9DgxocYLfpt!fACjVGGBeg2KwLCwDnWbhQCEN3qeUKiPg';
GRANT ALL PRIVILEGES ON rconfig.* TO 'rconfig'@'localhost';
FLUSH PRIVILEGES;
quit;
######################################################################################
curl -o- https://raw.githubusercontent.com/nvm-sh/nvm/v0.39.3/install.sh | bash
source ~/.bashrc
nvm ls-remote
nvm install node
######################################################################################
node -v
npm -v
######################################################################################
cd /var/www/html
git clone https://github.com/rconfig/rconfig.git
cd rconfig
cp .env.example .env
nano .env

APP_URL="https://rconfig.fbu.edu.tr"
DB_HOST=localhost
DB_PORT=3306
DB_DATABASE=rconfig
DB_USERNAME=rconfig
DB_PASSWORD=BKqD9DgxocYLfpt!fACjVGGBeg2KwLCwDnWbhQCEN3qeUKiPg
######################################################################################
export COMPOSER_ALLOW_SUPERUSER=1
composer self-update --2
yes | composer install --no-dev
######################################################################################
chmod +x setup_apache.sh
./setup_apache.sh
######################################################################################
nano /etc/apache2/sites-enabled/rconfig-vhost.conf
ServerName rconfig.fbu.edu.tr
ServerAlias rconfig.fbu.edu.tr

<VirtualHost *:443>
    ServerAdmin alisan.yildirim@fbu.edu.tr
    DocumentRoot "/var/www/html/rconfig/public"
    ServerName rconfig.fbu.edu.tr
    ServerAlias rconfig.fbu.edu.tr
    <Directory "/var/www/html/rconfig/public">
        Options Indexes FollowSymLinks
        AllowOverride All
        Require all granted
    </Directory>
    ErrorLog "/var/log/apache2/rconfig-error_log"
    CustomLog "/var/log/apache2/rconfig-access_log" combined
    SSLEngine on
    SSLCertificateFile /etc/certificate/fbu.edu.tr.crt
    SSLCertificateKeyFile /etc/certificate/fbu.edu.tr.key
    SSLCertificateChainFile /etc/certificate/ca.fbu.edu.tr.crt
</VirtualHost>
######################################################################################
mkdir /etc/certificate
(sertifikalar bu konuma kopyalanacak)
######################################################################################
sudo a2enmod ssl
######################################################################################
php artisan install
yes
######################################################################################
sudo systemctl restart apache2
######################################################################################
chown -R www-data /var/www/html/rconfig
######################################################################################
nano /etc/supervisor/conf.d/supervisord.conf

[program:horizon]
command = php /var/www/html/rconfig/artisan horizon
autostart = true
autorestart = true
stdout_logfile = /var/www/html/rconfig/logs/horizon.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile = /var/www/html/rconfig/logs/horizon.log
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0
######################################################################################
supervisorctl reread
supervisorctl update
supervisorctl start horizon
######################################################################################
Username: admin@domain.com
Password: admin
######################################################################################
http://rconfig.fbu.edu.tr
######################################################################################
Updating
cd /var/www/html/rconfig
git pull
php artisan migrate
php artisan rconfig:sync-tasks
composer install
npm install
npm run build
systemctl restart httpd
php artisan rconfig:clear-all
######################################################################################
